# systemd/can-simulation.service
[Unit]
Description=CAN Simulation for Development
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/thinkcity-dashboard-v3
# Setup vcan0
ExecStartPre=/sbin/modprobe vcan
ExecStartPre=/bin/sh -c 'ip link add dev vcan0 type vcan || true'
ExecStartPre=/sbin/ip link set up vcan0
# Start replay
ExecStart=/usr/bin/python3 /home/pi/thinkcity-dashboard-v3/tools/can_replay.py /home/pi/thinkcity-dashboard-v3/traces/entladen_91.log vcan0 50 --loop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
